# backup/systemd/pansy-phase-history.service
[Unit]
Description=PANSY transmitter phase history plot
After=pansy-backup-rsync.service
Wants=pansy-backup-rsync.service

[Service]
Type=simple
Environment=PANSY_RECEIVER_REPO=__PANSY_RECEIVER_REPO__
ExecStart=__PANSY_RECEIVER_REPO__/backup/scripts/pansy_phase_history_loop.sh
Restart=always
RestartSec=60

[Install]
WantedBy=default.target
