# receiver/systemd/pansy-uhd-rx.service
[Unit]
Description=PANSY UHD raw-voltage receiver
After=network.target

[Service]
Type=simple
Environment=PANSY_RECEIVER_REPO=__PANSY_RECEIVER_REPO__
Environment=PANSY_RECEIVER_CONFIG=%h/.config/pansy-receiver/pansy-receiver.env
WorkingDirectory=__PANSY_RECEIVER_REPO__
ExecStart=__PANSY_RECEIVER_REPO__/receiver/scripts/pansy_uhd_rx_service.sh
KillSignal=SIGINT
TimeoutStopSec=45
Restart=on-failure
RestartSec=10
StandardOutput=append:__PANSY_RECEIVER_REPO__/logs/pansy_uhd_rx.service.log
StandardError=append:__PANSY_RECEIVER_REPO__/logs/pansy_uhd_rx.service.log

[Install]
WantedBy=default.target
